# deploy/systemd/pwiki-vault-sync.service.example
[Unit]
Description=pwiki vault Git sync (host-side git pull)

[Service]
Type=oneshot
# A system unit runs as root by default, which usually does NOT have the ~/.ssh
# keys needed for an SSH Git remote. Set User=/Group= to the account that owns
# the vault checkout and can reach the remote. (Prefer the --user units in
# user-pwiki-vault-sync.*.example, which run as your login user automatically.)
#User=youruser
#Group=youruser
# Adjust these paths if the repository is not at /srv/newwiki. EnvironmentFile
# supplies PWIKI_GIT_HOST_DIR from the deployment .env.
WorkingDirectory=/srv/newwiki
EnvironmentFile=/srv/newwiki/.env
# Fast-forward pull on the host (NOT inside the container): the container image
# ships without an SSH client or keys. PWIKI_GIT_HOST_DIR must be an absolute
# path to the vault Git working tree. --ff-only refuses to run on a
# diverged/dirty tree instead of performing an automatic merge.
ExecStart=/bin/sh -lc 'git -C "$PWIKI_GIT_HOST_DIR" pull --ff-only'
