# deploy/systemd/user-pwiki-vault-sync.service.example
[Unit]
Description=pwiki vault Git sync (user)
After=default.target

[Service]
Type=oneshot
WorkingDirectory=/srv/pwiki
ExecStart=/usr/bin/docker compose exec -T pwiki sh -lc 'python -m pwiki.cli vault sync "$PWIKI_MARKDOWN_DIR"'
