# deploy/systemd/user-pwiki-vault-sync.timer.example
[Unit]
Description=Run pwiki vault Git sync every 5 minutes (user)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=pwiki-vault-sync.service

[Install]
WantedBy=timers.target
